# supervisord.conf
[supervisord]

[program:gateway]
command = python app.py
process_name = %(program_name)s_%(process_num)02d
numprocs = 7
autorestart = true